# systemd/GBackup.service
[Unit]
Description=Backup to Google Drive

[Service]
Type=simple
ExecStart=/root/backup_gdrive.sh

[Install]
WantedBy=multi-user.target
