# systemd/GBackup.timer
[Unit]
Description=Run GBackup every day at 01:00 and 13:00

[Timer]
OnCalendar=*-*-* 01,13:00:00

[Install]
WantedBy=timers.target
